# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autorestart=true